# lms/rootfs-ingress/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: LMS
# Runs the Nginx daemon
# ==============================================================================

# Wait for lms to become available
bashio::log.info "Waiting for lms before starting NGinx..."
bashio::net.wait_for 9000 localhost 900

bashio::log.info "Starting NGinx..."

exec nginx
