# lms/rootfs/etc/services.d/lms/run
#!/usr/bin/with-contenv bashio
set -e
set -u
set -o pipefail

readonly NAME=${HAA_NAME}

. /opt/${NAME}-support/lib/lib${NAME}.sh

function stop() {
  echo Signal
  kill $PID
}

function main() {
  # mmserver
  if [ "${LMS_mmserver:-''}" = "true" ];then
     bashio::log.info "Waiting for Music MagicIP"
     bashio::net.wait_for 10002 localhost 900
  fi
  bashio::log.info "Starting LMS..."
  # chnage buffering to line based for stdout sterr
  stdbuf -oL -eL squeezeboxserver ${LMS_OPTIONS} --user "${LMS_USER}" --group "${LMS_GROUP}" --prefsdir "${LMS_PREFS}" --logdir "${LMS_LOGDIR}" --cachedir "${LMS_CACHE}" </dev/null &
  PID=$!
  echo -n "$PID" > ${TMP:-"/tmp"}/.${NAME}_ready
  trap "stop" SIGTERM SIGINT
  trap "" SIGHUP

  wait "$PID"
  exit $? 
  ##FIXME
  #su ${LMS_USER} -s /bin/sh -c 'stdbuf -oL -eL squeezeboxserver ${LMS_OPTIONS} --user "${LMS_USER}" --group "${LMS_GROUP}" --prefsdir "${LMS_PREFS}" --logdir "${LMS_LOGDIR}" --cachedir "${LMS_CACHE}"'
  # --http-port? "${HTTP_PORT:-9000}"'
  #sleep 7d
}

main
# vim: ft=sh 
